# kubernetes/apps/home-automation/zwave/serial-forward/serial-forward.service
[Unit]
Description=Forward ttyAMA0 to TCP

[Service]
Type=simple
ExecStart=/usr/bin/serial-forward

[Install]
WantedBy=multi-user.target
